# ops/systemd/nullspace-ops.service
[Unit]
Description=Nullspace Ops Service
After=network.target

[Service]
Type=simple
User=nullspace
Group=nullspace
WorkingDirectory=/opt/nullspace/services/ops
EnvironmentFile=/etc/nullspace/ops.env
Environment=NODE_ENV=production
ExecStart=/usr/bin/node dist/server.js
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
