# ops/systemd/nullspace-executor.service
[Unit]
Description=Nullspace Executor Service
After=network.target

[Service]
Type=simple
User=nullspace
Group=nullspace
WorkingDirectory=/opt/nullspace
EnvironmentFile=/etc/nullspace/executor.env
ExecStart=/usr/local/bin/dev-executor --url ${EXECUTOR_URL} --identity ${EXECUTOR_IDENTITY} --block-interval-ms ${EXECUTOR_BLOCK_INTERVAL_MS}
Restart=on-failure
RestartSec=5
LimitNOFILE=100000

[Install]
WantedBy=multi-user.target
