# ops/systemd/nullspace-economy-snapshot.timer
[Unit]
Description=Run Nullspace economy snapshot every 5 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Unit=nullspace-economy-snapshot.service

[Install]
WantedBy=timers.target
