# ops/systemd/nullspace-node.service
[Unit]
Description=Nullspace Validator Node
After=network.target

[Service]
Type=simple
User=nullspace
Group=nullspace
WorkingDirectory=/opt/nullspace
EnvironmentFile=/etc/nullspace/node.env
Environment=NODE_ENV=production
ExecStart=/bin/sh -c '/usr/local/bin/nullspace-node --config "${NODE_CONFIG}" ${NODE_HOSTS:+--hosts "${NODE_HOSTS}"} ${NODE_PEERS:+--peers "${NODE_PEERS}"}'
Restart=on-failure
RestartSec=5
LimitNOFILE=100000

[Install]
WantedBy=multi-user.target
