# ops/systemd/docker/nullspace-auth-docker.service
[Unit]
Description=Nullspace Auth Service (Docker)
After=docker.service network.target
Requires=docker.service

[Service]
Type=simple
EnvironmentFile=/etc/nullspace/docker.env
EnvironmentFile=/etc/nullspace/auth.env
ExecStartPre=-/usr/bin/docker rm -f nullspace-auth
ExecStartPre=/usr/bin/docker pull ${IMAGE_REGISTRY}/nullspace-auth:${IMAGE_TAG}
ExecStart=/usr/bin/docker run --rm --name nullspace-auth --network host \
  --env-file /etc/nullspace/auth.env \
  -v /etc/nullspace:/etc/nullspace:ro \
  ${IMAGE_REGISTRY}/nullspace-auth:${IMAGE_TAG}
ExecStop=/usr/bin/docker stop nullspace-auth
Restart=on-failure
RestartSec=5
LimitNOFILE=100000

[Install]
WantedBy=multi-user.target
