# ops/systemd/docker/nullspace-gateway-docker.service
[Unit]
Description=Nullspace Gateway (Docker)
After=docker.service network.target
Requires=docker.service

[Service]
Type=simple
EnvironmentFile=/etc/nullspace/docker.env
EnvironmentFile=/etc/nullspace/gateway.env
ExecStartPre=-/usr/bin/docker rm -f nullspace-gateway
ExecStartPre=/usr/bin/docker pull ${IMAGE_REGISTRY}/nullspace-gateway:${IMAGE_TAG}
ExecStart=/usr/bin/docker run --rm --name nullspace-gateway --network host \
  --env-file /etc/nullspace/gateway.env \
  -v /var/lib/nullspace:/var/lib/nullspace \
  -v /etc/nullspace:/etc/nullspace:ro \
  ${IMAGE_REGISTRY}/nullspace-gateway:${IMAGE_TAG}
ExecStop=/usr/bin/docker stop nullspace-gateway
Restart=on-failure
RestartSec=5
LimitNOFILE=100000

[Install]
WantedBy=multi-user.target
