# ops/systemd/docker/nullspace-website-docker.service
[Unit]
Description=Nullspace Website (Docker)
After=docker.service network.target
Requires=docker.service

[Service]
Type=simple
EnvironmentFile=/etc/nullspace/docker.env
ExecStartPre=-/usr/bin/docker rm -f nullspace-website
ExecStartPre=/usr/bin/docker pull ${IMAGE_REGISTRY}/nullspace-website:${IMAGE_TAG}
ExecStart=/usr/bin/docker run --rm --name nullspace-website --network host \
  ${IMAGE_REGISTRY}/nullspace-website:${IMAGE_TAG}
ExecStop=/usr/bin/docker stop nullspace-website
Restart=on-failure
RestartSec=5
LimitNOFILE=100000

[Install]
WantedBy=multi-user.target
