# ops/systemd/docker/nullspace-simulator-docker.service
[Unit]
Description=Nullspace Simulator API (Docker)
After=docker.service network.target
Requires=docker.service

[Service]
Type=simple
EnvironmentFile=/etc/nullspace/docker.env
EnvironmentFile=/etc/nullspace/simulator.env
ExecStartPre=-/usr/bin/docker rm -f nullspace-simulator
ExecStartPre=/usr/bin/docker pull ${IMAGE_REGISTRY}/nullspace-simulator:${IMAGE_TAG}
ExecStart=/usr/bin/docker run --rm --name nullspace-simulator --network host \
  --env-file /etc/nullspace/simulator.env \
  --env NODE_ENV=production \
  -v /var/lib/nullspace:/var/lib/nullspace \
  -v /etc/nullspace:/etc/nullspace:ro \
  ${IMAGE_REGISTRY}/nullspace-simulator:${IMAGE_TAG} \
  --host ${SIMULATOR_HOST} --port ${SIMULATOR_PORT} --identity ${SIMULATOR_IDENTITY} ${SIMULATOR_ARGS}
ExecStop=/usr/bin/docker stop nullspace-simulator
Restart=on-failure
RestartSec=5
LimitNOFILE=100000

[Install]
WantedBy=multi-user.target
