# ops/systemd/docker/nullspace-ops-docker.service
[Unit]
Description=Nullspace Ops Service (Docker)
After=docker.service network.target
Requires=docker.service

[Service]
Type=simple
EnvironmentFile=/etc/nullspace/docker.env
EnvironmentFile=/etc/nullspace/ops.env
ExecStartPre=-/usr/bin/docker rm -f nullspace-ops
ExecStartPre=/usr/bin/docker pull ${IMAGE_REGISTRY}/nullspace-ops:${IMAGE_TAG}
ExecStart=/usr/bin/docker run --rm --name nullspace-ops --network host \
  --env-file /etc/nullspace/ops.env \
  -v /var/lib/nullspace:/var/lib/nullspace \
  -v /etc/nullspace:/etc/nullspace:ro \
  ${IMAGE_REGISTRY}/nullspace-ops:${IMAGE_TAG}
ExecStop=/usr/bin/docker stop nullspace-ops
Restart=on-failure
RestartSec=5
LimitNOFILE=100000

[Install]
WantedBy=multi-user.target
