# ops/systemd/docker/nullspace-node-docker.service
[Unit]
Description=Nullspace Validator Node (Docker)
After=docker.service network.target
Requires=docker.service

[Service]
Type=simple
EnvironmentFile=/etc/nullspace/docker.env
EnvironmentFile=/etc/nullspace/node.env
ExecStartPre=-/usr/bin/docker rm -f nullspace-node
ExecStartPre=/usr/bin/docker pull ${IMAGE_REGISTRY}/nullspace-node:${IMAGE_TAG}
ExecStart=/bin/sh -c '/usr/bin/docker run --rm --name nullspace-node --network host \
  --env-file /etc/nullspace/node.env \
  --env NODE_ENV=production \
  -v /var/lib/nullspace:/var/lib/nullspace \
  -v /etc/nullspace:/etc/nullspace:ro \
  ${IMAGE_REGISTRY}/nullspace-node:${IMAGE_TAG} \
  --config "${NODE_CONFIG}" ${NODE_HOSTS:+--hosts "${NODE_HOSTS}"} ${NODE_PEERS:+--peers "${NODE_PEERS}"}'
ExecStop=/usr/bin/docker stop nullspace-node
Restart=on-failure
RestartSec=5
LimitNOFILE=100000

[Install]
WantedBy=multi-user.target
