# ops/systemd/docker/nullspace-live-table-docker.service
[Unit]
Description=Nullspace Live Table Service (Docker)
After=docker.service network.target
Requires=docker.service

[Service]
Type=simple
EnvironmentFile=/etc/nullspace/docker.env
EnvironmentFile=/etc/nullspace/live-table.env
ExecStartPre=-/usr/bin/docker rm -f nullspace-live-table
ExecStartPre=/usr/bin/docker pull ${IMAGE_REGISTRY}/nullspace-live-table:${IMAGE_TAG}
ExecStart=/usr/bin/docker run --rm --name nullspace-live-table --network host \
  --env-file /etc/nullspace/live-table.env \
  -v /etc/nullspace:/etc/nullspace:ro \
  ${IMAGE_REGISTRY}/nullspace-live-table:${IMAGE_TAG}
ExecStop=/usr/bin/docker stop nullspace-live-table
Restart=on-failure
RestartSec=5
LimitNOFILE=100000

[Install]
WantedBy=multi-user.target
