# ops/systemd/docker/nullspace-executor-docker.service
[Unit]
Description=Nullspace Executor Service (Docker)
After=docker.service network.target
Requires=docker.service

[Service]
Type=simple
EnvironmentFile=/etc/nullspace/docker.env
EnvironmentFile=/etc/nullspace/executor.env
ExecStartPre=-/usr/bin/docker rm -f nullspace-executor
ExecStartPre=/usr/bin/docker pull ${IMAGE_REGISTRY}/nullspace-executor:${IMAGE_TAG}
ExecStart=/usr/bin/docker run --rm --name nullspace-executor --network host \
  --env-file /etc/nullspace/executor.env \
  -v /etc/nullspace:/etc/nullspace:ro \
  ${IMAGE_REGISTRY}/nullspace-executor:${IMAGE_TAG} \
  --url ${EXECUTOR_URL} --identity ${EXECUTOR_IDENTITY} --block-interval-ms ${EXECUTOR_BLOCK_INTERVAL_MS}
ExecStop=/usr/bin/docker stop nullspace-executor
Restart=on-failure
RestartSec=5
LimitNOFILE=100000

[Install]
WantedBy=multi-user.target
